# mkosi.extra/usr/lib/systemd/system/systemd-homed-firstboot.service.d/50-bananos.conf
# SPDX-License-Identifier: LGPL-2.1-or-later

# Upstream disables group prompting (--prompt-groups=no). Give the first
# interactive user the usual admin / device access groups and zsh instead.
# Missing groups are ignored by homectl.

[Service]
ExecStart=
ExecStart=/usr/lib/bananos/first-run
